# supervisor/pdf-api.conf
[program:pdf-api]
command=/usr/bin/python3 -m uvicorn api.main:app --host 127.0.0.1 --port 8001 --workers 2
directory=/www/wwwroot/pdf-to-text.digisource.vn
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=30
environment=HOME="/root",PYTHONPATH="/www/wwwroot/pdf-to-text.digisource.vn"
stdout_logfile=/var/log/supervisor/pdf-api.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/var/log/supervisor/pdf-api-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
